# templates/zpool-trim@.timer.j2
[Unit]
Description=Discard unused blocks in ZFS pool %I

[Timer]
OnCalendar={{ zfs_pool_trim_systemd_oncalendar }}
Persistent=true

[Install]
WantedBy=timers.target
