# templates/zpool-scrub@.timer.j2
[Unit]
Description=Scrub ZFS pool %i Service

[Timer]
OnCalendar={{ zfs_pool_scrub_systemd_oncalendar }}
Persistent=true

[Install]
WantedBy=timers.target